# iijp-api/docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log
[program:laravel-queue]
command=php /var/www/artisan queue:work --tries=3 --timeout=90
directory=/var/www
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel-queue.err.log
stdout_logfile=/var/log/laravel-queue.out.log
user=www-data